# systemd/translate.service
[Unit]
Description=TranslateGemma Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/translate-gemma
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="TRANSLATE_PORT=8080"
Environment="MODEL_DIR=/opt/translate-gemma/model"
ExecStart=/opt/translate-gemma/venv/bin/python /opt/translate-gemma/server.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
